# Ansible/populate_vm/templates/cri-dockerd.service.j2
[Unit]
Description=CRI Docker Daemon
After=network.target docker.service
Requires=docker.service

[Service]
ExecStart=/usr/local/bin/cri-dockerd

[Install]
WantedBy=multi-user.target
